# scripts/task.cron
0 0 * * * /bin/bash /home/naru/scripts/start_cron.sh
